# slow_talk09/ping1.service
[Unit]
Description= ping 1

[Service]
Type=simple
ExecStart=/root/temp/ping1.sh

[Install]
WantedBy=multi-user.target

